# docker/fts3-server/rootfs/etc/s6-overlay/s6-rc.d/fts-optimizer/run
#!/command/with-contenv sh
exec /usr/sbin/fts_optimizer -n
