# docker/fts3-server/rootfs/etc/s6-overlay/s6-rc.d/fts-server/run
#!/command/with-contenv sh
exec /usr/sbin/fts_server -n
