# docker/fts3-mon/rootfs/etc/s6-overlay/s6-rc.d/setupEnv/up
#!/command/with-contenv bash

exec s6-envdir /var/run/s6/container_environment /tmp/setup.sh
