# docker/fts3-rest/rootfs/etc/cron.d/httpd-cron
# min hr dom mon dow user  command
0 0,12 * * * root /command/s6-svc -1 /run/service/httpd >/dev/null 2>&1
